# aarch64/app/lib/tailscaled.service
[Unit]
Description=Tailscale node agent
Documentation=https://tailscale.com/kb/
Wants=network-pre.target
After=network-pre.target NetworkManager.service systemd-resolved.service

[Service]
Type=simple
ExecStart=/usr/local/packages/Tailscale_VPN/lib/tailscaled --state=/var/lib/tailscale/tailscaled.state
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target